# pytes.service.in
[Unit]
Description=PYTES BMS REST Service
After=network.target

[Service]
ExecStart=/usr/local/bin/pytes_daemon
WorkingDirectory=/usr/local/bin
StandardOutput=inherit
StandardError=inherit
Restart=always
User=$ENV{USER}

[Install]
WantedBy=multi-user.target
