# pytes.service.in
[Unit]
Description=PYTES BMS REST Service
After=network.target

[Service]
ExecStart=/opt/pytes/pytes_daemon
WorkingDirectory=/opt/pytes
StandardOutput=inherit
StandardError=inherit
Restart=always
User=$ENV{USER}

[Install]
WantedBy=multi-user.target
