# pytes.service.in
[Unit]
Description=PYTES BMS REST Service
After=network-online.target systemd-udev-settle.service
Requires=systemd-udev-settle.service
Wants=network-online.target

[Service]
ExecStart=/opt/pytes/pytes_daemon
WorkingDirectory=/opt/pytes
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=5
StartLimitIntervalSec=0
StartLimitBurst=0
User=$ENV{USER}

[Install]
WantedBy=multi-user.target
